# supervisord-gateone.conf
[program:gateone]
command=/start-gateone.sh
numprocs=1
autostart=true
autorestart=true