# supervisord-cow.conf
[program:cow]
command=/root/cow
numprocs=1
autostart=true
autorestart=true