# superstart/supervisord-Mammoth_titan.conf
[program:Mammoth_titan]
command=php /root/myenv/cloud/coding/Mammoth_titan/socketio_server.php start
numprocs=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
