# superstart/supervisord-Mammoth_tty.conf
[program:Mammoth_tty]
command=php /root/myenv/cloud/coding/Mammoth_tty/start.php start
numprocs=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
