# superstart/supervisord-logstash.conf
[program:logstash]
command=/opt/logstash/bin/logstash -f /etc/logstash/conf.d/logstash.conf
autostart=false
autorestart=true