# superstart/supervisord-Mammoth_client.conf
[program:Mammoth_client]
command=python app.py
numprocs=1
autostart=true
autorestart=true
directory=/root/myenv/cloud/coding/Mammoth_client
